# roles/forgejo/files/forgejo.service
[Unit]
Description=Forgejo self-hosted Git service
After=network.target mariadb.service

[Service]
Type=simple
User=git
Group=git
WorkingDirectory=/var/lib/forgejo
ExecStart=/usr/local/bin/forgejo web --config /etc/forgejo/app.ini --work-path /opt/forgejo
Environment=USER=git
Environment=HOME=/home/git
Environment=FORGEJO_WORK_DIR=/var/lib/forgejo
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
